# debian/rootfs/etc/s6/s6-init/run
#!/bin/bash

# shellcheck disable=SC1091

set -o errexit
set -o nounset
set -o pipefail

[ -n "${DEBUG:+1}" ] && set -x

# Load libraries
. /opt/easysoft/scripts/liblog.sh
. /opt/easysoft/scripts/libmysql.sh
. /opt/easysoft/scripts/libpersistence.sh

# Load Global/Apache/PHP environment variables
. /etc/s6/s6-init/envs

#==================================#
# Prepare persistence directories. #
#==================================#
info "Prepare persistence directories."
for pdir in "${PHP_PERSISTENCE_LIST[@]}"
do
    # Ensure a directory exists and,is owned by the given user
    ensure_dir_exists "$pdir" "www-data"
done

# Make soft link for persistent directory
for ldir in "${PERSISTENCE_LINK_LIST[@]}"
do
    ensure_dir_exists "$(dirname $ldir)" "www-data"
    move_then_link "$ldir" "${ldir/data/apps}" "www-data"
done

#=====================#
#   Prepare PHP       #
#=====================#
info "Render php.ini with environment variables."
if [ "$OS_ARCH" == "arm64" ];then
    info "render template php.ini for arm64."
    PHP_INI="/usr/local/etc/php/php.ini"
    /usr/bin/render-template ${PHP_INI}.tpl > ${PHP_INI}
else
    info "render template php.ini for amd64"
    PHP_CLI_INI="/etc/php/7.4/cli/php.ini"
    PHP_APACHE_INI="/etc/php/7.4/apache2/php.ini"
    /usr/bin/render-template /etc/php/7.4/php.ini.tpl > ${PHP_CLI_INI}
    /usr/bin/render-template /etc/php/7.4/php.ini.tpl > ${PHP_APACHE_INI}
fi

# enable redis extension
if [ "$PHP_SESSION_TYPE" == "redis" ];then
    info "Enable redis extension for session."
    phpenmod redis
fi

#=====================#
#   Prepare Apache    #
#=====================#
PORTS_CFG="/etc/apache2/ports.conf"
SITE_CFG="/etc/apache2/sites-enabled/site.conf"
info "Render apache sites config with envionment variables."
/usr/bin/render-template ${PORTS_CFG}.tpl > $PORTS_CFG
/usr/bin/render-template ${SITE_CFG}.tpl > $SITE_CFG

#======================#
# Check MySQL service  #
#======================#
# Check and waiting mysql to be ready
wait_for_mysql || exit 1

# Initialize database
mysql_init_db "$MYSQL_DB"


#=====Set Zentao Config=====
# Disable checkClient fix QuickOn DevOps
cfg_num=$(grep -n checkClient /apps/zentao/config/config.php | awk -F ':'  '{print $1}')
sed -i "$cfg_num s/true/false/" /apps/zentao/config/config.php
